# deploy/supervisor/reverb.example.conf
[program:audit-platform-reverb]
command=php /var/www/audit-platform/artisan reverb:start --host=0.0.0.0 --port=8080
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/audit-platform/storage/logs/reverb.log
